# cb-yara-connector.service
[Unit]
Description=Carbon Black Response Yara Connector
After=syslog.target network.target

[Service]
Environment=C_FORCE_ROOT=1
Type=forking
ExecStart=/usr/share/cb/integrations/cb-yara-connector/yaraconnector --config-file /etc/cb/integrations/cb-yara-connector/yaraconnector.conf --run-forever True --log-file /var/log/cb/integrations/cb-yara-connector/yaraconnector.log --output-file /var/cb/data/cb-yara-connector/feed.json
WorkingDirectory=/usr/share/cb/integrations/cb-yara-connector
KillSignal=SIGQUIT


[Install]
WantedBy=multi-user.target
